# citofono-voip.service
[Unit]
Description=Sistema Citofono VoIP
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/citofono-voip
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /opt/citofono-voip/citofono-voip.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Limiti risorse
MemoryLimit=128M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
